# templates/unicorn.init.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:             unicorn
# Required-Start:       $all
# Required-Stop:        $all
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Unicorn master server
### END INIT INFO

NAME=unicorn
DESC="unicorn"
CONFIG="<%= "#{configuration[:deploy_to]}/shared/config/unicorn.rb" %>"
PIDFILE="<%= "#{configuration[:deploy_to]}/shared/pids/unicorn.pid" %>"
BINARY="<%= "#{configuration[:deploy_to]}/shared/binary/#{configuration[:unicorn][:binary]}" %>"

set -u
set -e
# Example init script, this can be used with nginx, too,
# since nginx and unicorn accept the same signals

# Make sure the binary and the config file are present before proceeding
test -x $BINARY || exit 0
test -f $CONFIG || exit 0

# Feel free to change any of the following variables for your app:
APP_ROOT="<%= "#{configuration[:deploy_to]}/current" %>"
ENV="<%= configuration[:unicorn][:rails_env] || 'production' %>"
CMD="$BINARY -D -E $ENV -c $CONFIG"

old_pid="$PIDFILE.oldbin"

cd $APP_ROOT || exit 1

sig () {
	test -s "$PIDFILE" && kill -$1 `cat $PIDFILE`
}

oldsig () {
	test -s $old_pid && kill -$1 `cat $old_pid`
}

case $1 in
start)
  echo -n "Starting $DESC: "
	sig 0 && echo >&2 "Already running" && exit 0
	$CMD
	echo "$NAME."
	;;
stop)
  echo -n "Stopping $DESC: "
	sig QUIT && exit 0
	echo "$NAME."
	;;
force-stop)
echo -n "Force-stopping $DESC: "
	sig TERM && exit 0
	echo "$NAME."
	;;
restart|reload)
  echo -n "Restarting $DESC: "
	sig HUP && echo "$NAME." && exit 0
	echo -n >&2 "Couldn't reload, starting '$CMD' instead: "
	$CMD
	echo "$NAME."
	;;
upgrade)
  echo -n "Upgrading $DESC: "
	sig USR2 && echo "$NAME." && exit 0
	echo -n >&2 "Couldn't upgrade, starting '$CMD' instead: "
	$CMD
	echo "$NAME."	
	;;
rotate)
  sig USR1 && echo rotated logs OK && exit 0
  echo >&2 "Couldn't rotate logs" && exit 1
  ;;
*)
	echo >&2 "Usage: $0 <start|stop|restart|upgrade|rotate|force-stop>"
	exit 1
	;;
esac